# task-mgr.service
[Unit]
Description=Rust Task Manager Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/task-mgr
ExecStart=/opt/task-mgr/server/task-mgr-server
Restart=always
RestartSec=3
Environment=RUST_LOG=info
Environment=PORT=3000
Environment=WEB_ROOT=/opt/task-mgr/web
Environment=DATABASE_URL=sqlite:///opt/task-mgr/tasks.db

[Install]
WantedBy=multi-user.target
